# owncloud/rootfs/etc/s6-overlay/s6-rc.d/php-fpm/run
#!/usr/bin/env bashio
# ==============================================================================
# Home Assistant Add-on: ownCloud  
# Runs PHP-FPM
# ==============================================================================

bashio::log.info "Starting PHP-FPM..."

# Create PHP-FPM socket directory
mkdir -p /var/run/php

# Start PHP-FPM
exec php-fpm82 --nodaemonize --fpm-config /etc/php82/php-fpm.conf